# src/etc/profile.d/r4pi.sh
# Added for the r4pi toolkit - R for the Raspberry Pi
# See https://r4pi.org for more information
export PATH=/opt/r4pi/bin/R:${PATH}
